# onboard/video/src/supervisord.conf
[supervisord]
nodaemon=true

[program:streamer]
command=/usr/local/bin/streamer.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/streamer.err.log
stdout_logfile=/var/log/streamer.out.log

